# stage2/04-pirate-radio/files/etc/systemd/system/wpa_supplicant@wlan0.service.d/override.conf
[Unit]
BindsTo=accesspoint@%i.service
After=accesspoint@%i.service
